# config-files/etc/profile.d/composer.sh
#!/bin/bash

export COMPOSER_HOME=/usr/share/composer
